# files/macos/Library/LaunchDaemons/time-machine-metadata-exporter.plist
<?xml version="1.0" encoding="UTF-8"?>
  <!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
  <plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.djrtechconsulting.time_machine_metadata_exporter</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/local/sbin/time-machine-metadata-exporter</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/time_machine_metadata_exporter.err.log</string>
    <key>StandardOutPath</key>
    <string>/usr/local/var/log/time_machine_metadata_exporter.log</string>
  </dict>
</plist>
