# templates/macos/Library/LaunchDaemons/com.djrtechconsulting.node_exporter_shell.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
  <!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
  <plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.djrtechconsulting.node_exporter_shell</string>
    <key>ProgramArguments</key>
    <array>
      <string>sh</string>
      <string>-c</string>
      <string>/usr/local/sbin/node-exporter-shell</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>{{ homebrew_prefix }}/bin:{{ ansible_env.PATH }}</string>
    </dict>
    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/node_exporter_shell.err.log</string>
    <key>StandardOutPath</key>
    <string>/usr/local/var/log/node_exporter_shell.log</string>
  </dict>
</plist>
