# templates/Library/LaunchDaemons/com.wireguard.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd";>
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.wireguard.{{ wg_interface }}</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/sbin/wireguard_runner</string>
  </array>
  <key>OnDemand</key>
  <false/>
  <key>RunAtLoad</key>
  <true/>
  <key>TimeOut</key>
  <integer>90</integer>
  <key>KeepAlive</key>
  <false/>
  <key>AbandonProcessGroup</key>
  <true/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>{{ homebrew_prefix }}/bin:{{ homebrew_prefix }}/sbin:/usr/local/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>wg_interface</key>
    <string>{{ wg_interface }}</string>
  </dict>
  <key>StandardErrorPath</key>
  <string>/var/log/wireguard.{{ wg_interface }}.err</string>
  <key>StandardOutPath</key>
  <string>/var/log/wireguard.{{ wg_interface }}.out</string>
</dict>
</plist>
